# deploy/sector-tracker-stops.service
[Unit]
Description=Enforce declared Fringe stops against live marks
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/ds/.local/bin/fringe_stop_monitor.py
