# deploy/sector-tracker-uploader.path
[Unit]
Description=Watch the synced Hermes vault for new reports

[Path]
PathModified=/opt/sector-tracker/vault
Unit=sector-tracker-uploader.service

[Install]
WantedBy=multi-user.target
