# deploy/sector-tracker-uploader.path
[Unit]
Description=Watch the Hermes vault for new cron reports

[Path]
PathModified=/home/ds/hermes-research
Unit=sector-tracker-uploader.service

[Install]
WantedBy=default.target
